# deploy/templates/service-graph-refresh@.timer.template
[Unit]
Description=Scheduled Service Graph refresh for %i

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=20m
Persistent=true
Unit=service-graph-refresh@%i.service

[Install]
WantedBy=timers.target
